# automatic-timezone/config/etc/NetworkManager/dispatcher.d/09-timezone
#!/bin/sh

case "$2" in
    connectivity-change)
        tzupdate
    ;;
esac
